# rootfs/etc/systemd/system/clean.timer
[Unit]
Description=Run clean.service daily

[Timer]
OnCalendar=daily
AccuracySec=12h
Persistent=true

[Install]
WantedBy=timers.target
