# rootfs/etc/systemd/system/nvclock.service
[Unit]
Description=Set NVIDIA overclocks

[Service]
ExecStart=/etc/systemd/system/nvclock.bin

[Install]
WantedBy=multi-user.target
