# overclock/overclock.service
[Unit]
Description=Set NVIDIA overclocks

[Service]
ExecStart=/bin/overclock

[Install]
WantedBy=multi-user.target
