# overclock/overclock.service
[Unit]
Description=Set NVIDIA overclocks

[Service]
ExecStart=/usr/bin/overclock

[Install]
WantedBy=multi-user.target
