# system/systemd/system/overclock.service
[Unit]
Description=Set NVIDIA overclocks

[Service]
ExecStart=/etc/systemd/system/overclock.bin

[Install]
WantedBy=multi-user.target
