# etc/systemd/system/overclock.service
[Unit]
Description=Set NVIDIA overclocks

[Service]
ExecStart=/usr/local/bin/overclock

[Install]
WantedBy=multi-user.target
